# systemd/agent-renew.timer
[Unit]
Description=Refresh scanner /loop cron every 6 days (claude /loop crons auto-expire at 7 days)

[Timer]
# Fire 5 minutes after boot in case the in-memory cron was lost while the box was off.
OnBootSec=5min
# Then every 6 days after the last renewal run — stays comfortably below the 7-day TTL.
OnUnitActiveSec=6d
# Add some randomization so coordinated clusters don't all refresh at once.
RandomizedDelaySec=15min
Unit=agent-renew.service

[Install]
WantedBy=timers.target
