# systemd/agent-renew.service
[Unit]
Description=Kill scanner tmux session so supervised-agent@scanner.service re-registers the /loop cron before the 7-day TTL expires
Requisite=supervised-agent@scanner.service
After=supervised-agent@scanner.service

[Service]
Type=oneshot
User=dev
Group=dev
# Skip if agent is paused — no point cycling a paused session.
ExecStart=/bin/sh -c 'if [ -f /var/run/kick-governor/paused_scanner ]; then echo "scanner is paused — skipping renew"; exit 0; fi; /usr/bin/tmux kill-session -t scanner 2>/dev/null || true'
SuccessExitStatus=0 1
