# systemd/kick-governor.timer
[Unit]
Description=Run kick-governor every 15 minutes to adaptively kick supervised agents
Requires=kick-governor.service

[Timer]
# Aligns to wall-clock :00/:15/:30/:45 so timing is predictable.
OnCalendar=*:00,15,30,45
AccuracySec=30s
Unit=kick-governor.service

[Install]
WantedBy=timers.target
