# systemd/kick-outreach.service
[Unit]
Description=Kick outreach agent to run an organic search/traffic pass
After=network-online.target

[Service]
Type=oneshot
User=dev
Environment=TMUX_BIN=tmux
ExecStart=/usr/local/bin/kick-agents.sh outreach
StandardOutput=journal
StandardError=journal
