# systemd/kick-governor.service
[Unit]
Description=Adaptive kick governor — adjusts agent cadences based on issue/PR queue depth
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=dev
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
EnvironmentFile=-/etc/hive/governor.env
ExecStart=/usr/local/bin/kick-governor.sh
StandardOutput=journal
StandardError=journal
