# systemd/supervised-agent-healthcheck.service
[Unit]
Description=Supervised-agent healthcheck (detect heartbeat stall, respawn, notify)
Requisite=supervised-agent.service
After=supervised-agent.service

[Service]
Type=oneshot
User=__AGENT_USER__
Group=__AGENT_USER__
EnvironmentFile=/etc/supervised-agent/agent.env
ExecStart=/usr/local/bin/agent-healthcheck.sh
