# systemd/supervised-agent-renew@.timer
[Unit]
Description=Refresh /loop cron for supervised-agent instance %i every 6 days

[Timer]
OnBootSec=5min
OnUnitActiveSec=6d
RandomizedDelaySec=15min
Unit=supervised-agent-renew@%i.service

[Install]
WantedBy=timers.target
